# model/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[inet_http_server]         ; inet (TCP) server disabled by default
port=*:9001        ; (ip_address:port specifier, *:port for all iface)
username=123              ; (default is no username (open server))
password=123               ; (default is no password (open server))

[supervisord]
logfile=/root/project/news/logs/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/root/project/news/logs/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///root/project/news/logs/supervisor.sock ; use a unix:// URL  for a unix socket

[program:news_gunicorn]
directory=/root/project/news/
command=gunicorn -c /etc/gunicorn.conf NewsStatus.wsgi:application
numprocs=1
;process_name=999
autorestart=true
autostart=true
startsecs=30
;startretries=3
stdout_logfile=/root/project/news/logs/gunicorn_out.log
stderr_logfile=/root/project/news/logs/gunicorn_err.log



[program:celery]
command=/root/project/news/cron.sh start
autostart=true
autorestart=true
priority=10
stdout_events_enabled=true
stderr_events_enabled=true
